# app_eval_1.service
[Unit]
Description=App Eval 1
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/java -jar /artifact/app_eval_1-0.0.1-SNAPSHOT.jar
[Install]
WantedBy=multi-user.target